# crontab
@reboot  /Users/av/scripts/osx/keyboard.sh
